# upi-fraud-advanced.service
[Unit]
Description=UPI Fraud Detection Advanced System
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/upi-fraud-system
ExecStart=/usr/bin/python3 advanced_quick_start.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
